# payload/etc/systemd/system/systemd-pcrlogin@.service.d/10-spawn-arch-disable.conf
[Unit]
ConditionPathExists=/var/lib/spawn-arch/enable-pcrlogin
