# payload/etc/systemd/system/spawn-arch-btrfs-scrub.service
[Unit]
Description=Scrub the root Btrfs filesystem
ConditionPathIsMountPoint=/

[Service]
Type=oneshot
ExecStart=/usr/bin/btrfs scrub start -B /
