# payload/etc/systemd/system/docker.service.d/10-spawn-arch-ordering.conf
[Unit]
Wants=network-online.target firewalld.service
After=network-online.target firewalld.service
